# systemd/natgeo-wallpaper.service
[Unit]
Description=Set NatGeo Photo of the Day as desktop wallpaper
After=network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/national-geographic-api
ExecStart=/usr/bin/npm run wallpaper
# Give GNOME's settings daemon access to the session D-Bus
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
