# systemd/natgeo-wallpaper.timer
[Unit]
Description=Daily NatGeo wallpaper update
Requires=natgeo-wallpaper.service

[Timer]
# Run once a day at 08:00 (after NatGeo typically publishes the new photo)
OnCalendar=*-*-* 08:00:00
# If the machine was off at 08:00, run as soon as it boots
Persistent=true
Unit=natgeo-wallpaper.service

[Install]
WantedBy=timers.target
